# drone-updater.service
[Unit]
Description=nRF OTA Firmware Auto Updater Service for Drones
# Wait for network and bluetooth
After=network.target bluetooth.target
# Wait for downloader service if it exists/is enabled
After=firmware-downloader.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/drone_updater
ExecStartPre=+/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/python3 /opt/drone_updater/drone_updater.py
Restart=always
RestartSec=5
# Output standard logs to journal as well
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
